# tasks/icon-cache-event-repair.xml
<?xml version="1.0" encoding="UTF-16"?>
<!--
  icon-cache-event-repair.xml
  Layer A — Event-Triggered Repair (Reference Copy)

  NOTE: In v2.0.0+, Register-Tasks.ps1 generates the task XML inline
  and registers it directly via schtasks.exe. This file is retained
  as a readable reference showing the full task definition for Layer A.

  Triggers:
    - Event ID 1000 (Application Error) — explorer.exe hard crash
    - Event ID 1002 (Application Hang)  — explorer.exe stopped responding
    - Event ID 107  (Kernel-Power)      — system resumed from sleep

  Naming Policy: naming-conventions-policy-v3.2.0 — Style B (kebab-case)
-->
<Task version="1.4" xmlns="http://schemas.microsoft.com/windows/2004/02/mit/task">
  <RegistrationInfo>
    <Description>icon-cache-self-healing v2.0.0 - Layer A. Repairs icon cache after explorer.exe crash (Event 1000), hang (Event 1002), or sleep resume (Event 107).</Description>
    <URI>\IconCache\EventRepair</URI>
  </RegistrationInfo>
  <Triggers>
    <EventTrigger>
      <Enabled>true</Enabled>
      <Subscription><![CDATA[<QueryList><Query Id="0" Path="Application"><Select Path="Application">*[System[Provider[@Name='Application Error'] and EventID=1000]] and *[EventData[Data[@Name='param1'] and (Data='explorer.exe')]]</Select></Query></QueryList>]]></Subscription>
      <Delay>PT60S</Delay>
    </EventTrigger>
    <EventTrigger>
      <Enabled>true</Enabled>
      <Subscription><![CDATA[<QueryList><Query Id="0" Path="Application"><Select Path="Application">*[System[Provider[@Name='Application Hang'] and EventID=1002]] and *[EventData[Data[@Name='param1'] and (Data='explorer.exe')]]</Select></Query></QueryList>]]></Subscription>
      <Delay>PT90S</Delay>
    </EventTrigger>
    <EventTrigger>
      <Enabled>true</Enabled>
      <Subscription><![CDATA[<QueryList><Query Id="0" Path="System"><Select Path="System">*[System[Provider[@Name='Microsoft-Windows-Kernel-Power'] and EventID=107]]</Select></Query></QueryList>]]></Subscription>
    </EventTrigger>
  </Triggers>
  <Principals>
    <Principal id="Author">
      <LogonType>InteractiveToken</LogonType>
      <RunLevel>LeastPrivilege</RunLevel>
    </Principal>
  </Principals>
  <Settings>
    <MultipleInstancesPolicy>IgnoreNew</MultipleInstancesPolicy>
    <DisallowStartIfOnBatteries>false</DisallowStartIfOnBatteries>
    <StopIfGoingOnBatteries>false</StopIfGoingOnBatteries>
    <Hidden>true</Hidden>
    <ExecutionTimeLimit>PT5M</ExecutionTimeLimit>
    <Priority>7</Priority>
  </Settings>
  <Actions>
    <Exec>
      <Command>pwsh.exe</Command>
      <Arguments>-WindowStyle Hidden -NonInteractive -ExecutionPolicy Bypass -File "scripts\Repair-IconCache.ps1"</Arguments>
    </Exec>
  </Actions>
</Task>
